# portainer/rootfs/etc/s6-overlay/s6-rc.d/init-nginx/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Community App: Portainer
# Configures Nginx for use with the Portainer UI
# ==============================================================================

# Generate Ingress configuration
bashio::var.json \
    interface "$(bashio::app.ip_address)" \
    port "^$(bashio::app.ingress_port)" \
    | tempio \
        -template /etc/nginx/templates/ingress.gtpl \
        -out /etc/nginx/servers/ingress.conf
